# config_files/crontab
* * * * * sha1sum -c /tmp/checksums | grep -v "OK" >> /var/log/checksums.log
